# examples/dtn-discovery.service
[Unit]
Description=DTN Neighbor Discovery Daemon
After=network.target ionwd.service dtnex.service
Wants=dtnex.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/dtn/dtn-discovery
ExecStart=/usr/bin/python3 %h/dtn/dtn-discovery/discovery.py
Restart=on-failure
RestartSec=30
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
